# backend/deploy/securetarget-backend.service
[Unit]
Description=SecureTarget backend API
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/SecureTarget
EnvironmentFile=/etc/securetarget.env
ExecStart=/usr/bin/npm run start --workspace @securetarget/backend
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
